# modules/systemd/earnbox-api.service
[Unit]
Description=EarnBox Backend API
After=network.target docker.service
Requires=docker.service

[Service]
Type=simple
WorkingDirectory=/home/pi/EarnBox
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

# Ensure logs go to journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
