# modules/systemd/earnbox-reset.service
[Unit]
Description=EarnBox Daily Reset Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/home/pi/EarnBox
ExecStart=/usr/bin/docker compose down
ExecStartPost=/usr/bin/docker compose up -d

# Log to journalctl
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
